# workshop/supervisor/explorer.conf
[program:quorum]
directory=/home/eduk8s/resources/explorer
command=npm run dev 
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true